# dawarich/rootfs/etc/s6-overlay/s6-rc.d/init-dawarich/up
/etc/s6-overlay/scripts/init-dawarich.sh
