# dawarich/rootfs/etc/s6-overlay/s6-rc.d/svc-ingress/run
#!/usr/bin/with-contenv bash
# Wait for the generated nginx config (written by init-dawarich.sh)
while [ ! -f /etc/nginx/nginx.conf ]; do
  sleep 1
done
exec nginx -c /etc/nginx/nginx.conf
