# dawarich/rootfs/etc/s6-overlay/s6-rc.d/svc-sidekiq/finish
#!/usr/bin/with-contenv bash
# Give Sidekiq time to finish in-progress jobs before s6 sends SIGKILL
# Sidekiq traps SIGTERM and stops picking new jobs, then waits for current ones
sleep 10
