# dawarich/rootfs/etc/s6-overlay/s6-rc.d/svc-postgres/finish
#!/usr/bin/with-contenv bash
# Gracefully stop PostgreSQL on container shutdown
if [ -f /data/postgres/postmaster.pid ]; then
  su - postgres -c "PATH=/usr/lib/postgresql/17/bin:\$PATH pg_ctl stop -D /data/postgres -m fast"
fi
